# systemd/homebackup-snapshot.timer
[Unit]
Description=Daily snapshot of backup drive at 3:00 AM (after backup window)

[Timer]
# Backups run 2:00–2:15 AM. Snapshot fires at 3:00 AM once all machines are done.
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=homebackup-snapshot.service

[Install]
WantedBy=timers.target
